# scripts/filebrowser.service
[Unit]
Description=File Browser (file manager behind /files)
After=network.target

[Service]
Type=simple
# Bound to 127.0.0.1 only (see /etc/filebrowser/filebrowser.json). The liveoaks
# app reverse-proxies /files to it and injects X-Forwarded-User; nothing else
# may reach it, or the proxy-auth header could be spoofed.
ExecStart=/usr/local/bin/filebrowser -c /etc/filebrowser/filebrowser.json
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
